# app/src/main/deploy/package/linux/postinst
#!/bin/sh
# postinst script for WT4
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package

case "$1" in
    configure)
        echo Adding shortcut to the menu


        xdg-desktop-menu install --novendor /opt/WT4/WT4.desktop


        if [ "false" = "true" ]; then
            echo Installing daemon
            cp /opt/WT4/wt4.init /etc/init.d/wt4

            if [ -x "/etc/init.d/wt4" ]; then
                update-rc.d wt4 defaults

                if [ "false" = "true" ]; then
                    if which invoke-rc.d >/dev/null 2>&1; then
                        invoke-rc.d wt4 start
                    else
                        /etc/init.d/wt4 start
                    fi
                fi
	        fi

        fi

	#PERM_USER=$(ps -o user= -p $$ | awk '{print $1}')
	#echo "Changing user permissions as: $PERM_USER"
	#sudo chown -R ${PERM_USER}:${PERM_USER} /opt/WT4

	echo "Adding permissions for update folder"
	sudo chmod -R 0777 /opt/WT4/app

	JRE_DIR="/opt/WT4/runtime"
	if [ -d "$JRE_DIR" ]; then
		echo "Using bundled JRE"
	else
		echo "Adding local classpath to JRE"
		HACKY_JAVA_HOME=$(dirname $(dirname $(readlink -e /usr/bin/javac)))
		if [ ! -z $JAVA_HOME ] 
		then
			echo "Using JRE from: $JAVA_HOME"
			sudo sed -i "s|app.runtime=.*|app.runtime=$JAVA_HOME|g" /opt/WT4/app/WT4.cfg
		elif [ ! -z $HACKY_JAVA_HOME ]
		then
			echo "Using JRE from: $HACKY_JAVA_HOME"
			sudo sed -i "s|app.runtime=.*|app.runtime=$HACKY_JAVA_HOME|g" /opt/WT4/app/WT4.cfg
		else
			echo "Warning: JAVA_HOME not set!"
		fi
	fi
	;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
